# challenges/ctf-forensics-dfir-elk/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB

[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
environment=ES_JAVA_OPTS="-Xms512m -Xmx512m"
stdout_logfile=/var/log/elasticsearch.log
stderr_logfile=/var/log/elasticsearch.log
autorestart=true
priority=10
startretries=3

[program:kibana]
command=/usr/share/kibana/bin/kibana
user=kibana
stdout_logfile=/var/log/kibana.log
stderr_logfile=/var/log/kibana.log
autorestart=true
priority=20
startsecs=30

[program:inject-logs]
command=/usr/bin/python3 /opt/inject_logs.py
user=root
stdout_logfile=/var/log/inject-logs.log
stderr_logfile=/var/log/inject-logs.log
autorestart=false
startsecs=0
priority=30

[program:kibana-init]
command=/opt/import_kibana.sh
stdout_logfile=/var/log/kibana-init.log
stderr_logfile=/var/log/kibana-init.log
autorestart=false
startsecs=5
priority=40
